# cloudrun/tags/cloudrun-0.9/upstart/bootstrap.conf
#upstart-job
description "Schedule execution of bootstrap.sh"
task
start on stopped cloud-final
stop on runlevel [016]
console output
pre-start script
  echo 'pre-start bootstrap.conf...'
end script
exec start-stop-daemon --start --chdir /home/ubuntu --chuid ubuntu --exec /home/ubuntu/bootstrap.sh
post-stop script
  echo 'post-stop bootstrap.conf'
end script


